# platform/systemd/goodlove-auto-deploy.timer
[Unit]
Description=Check GitHub main for Goodlove PWA updates

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
RandomizedDelaySec=15s
Persistent=true
Unit=goodlove-auto-deploy.service

[Install]
WantedBy=timers.target
